# config/autostart/antimicrox.appimage.desktop
[Desktop Entry]
Exec=/home/taylor/AppImages/antimicrox.appimage
Icon=
Name=antimicrox.appimage
Path=
Terminal=False
Type=Application
